# .config/systemd/user/schedule-notes.service
[Unit]
Description=Notify based on journal notes

[Service]
Type=simple
ExecStart=/home/gustavo/.bin/notes

[Install]
WantedBy=default.target
