# systemd/notify-journal-todo.service
[Unit]
Description=Notify todo

[Service]
ExecStart=/home/gustavo/dotfiles-public/scripts/journal-notify.sh
User=gustavo
