# systemd/sync-notes.service
[Unit]
Description=Sync notes with GitHub

[Service]
ExecStart=/home/gustavo/dotfiles-public/scripts/sync-notes.sh
User=gustavo
